# shell/system-backup.service
[Unit]
Description=System Settings Backup

[Service]
Type=oneshot
ExecStart=/home/svartfax/.shell/system_backup.sh

